# resources/etc/supervisor/conf.d/hbase.conf
[program:hbase]
command=/app/start-hbase.sh
user=worker
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=JAVA_HOME="/usr/lib/jvm/java-8-oracle"
