# resources/etc/supervisor/conf.d/kafka.conf
[program:kafka]
command=/app/start-kafka.sh
user=worker
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

